# float/roles/base/templates/profile.d/logcat.sh.j2
export ELASTICSEARCH_URL="http://log-collector.{{ domain }}:9200"
