# float/roles/float-infra-nginx/files/acme-reload-nginx.cron
45 3 * * * root /usr/local/bin/acme-reload-nginx
