# float/roles/float-base/files/float-lockdown.target
[Unit]
Description=float has been locked down
Requires=multi-user.target
After=multi-user.target

[Install]
WantedBy=multi-user.target
